# infra/launchd/ai.jarvis.daily-restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.jarvis.daily-restart</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ramsbaby/jarvis/runtime/scripts/bot-daily-restart.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/daily-restart.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/daily-restart.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
