# infra/launchd/ai.jarvis.report-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.jarvis.report-monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ramsbaby/jarvis/runtime/scripts/report-generate.sh</string>
        <string>monthly</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/ramsbaby</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/report-monthly.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/report-monthly.err.log</string>

    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
