# infra/launchd/ai.jarvis.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.jarvis.discord-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ramsbaby/jarvis/runtime/scripts/bot-preflight.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ramsbaby/jarvis/runtime/discord</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ramsbaby</string>
        <key>BOT_HOME</key>
        <string>/Users/ramsbaby/jarvis/runtime</string>
        <key>JARVIS_HOME</key>
        <string>/Users/ramsbaby/jarvis</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/Users/ramsbaby/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/discord-bot.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ramsbaby/jarvis/runtime/logs/discord-bot.err.log</string>
</dict>
</plist>
